# systemd/clone-pi.service
[Unit]
Description=Run clone-pi.sh

[Service]
Type=oneshot
ExecStart=/home/jiblet/bin/clone-pi.sh --force
