# systemd/discord-startup@.service
[Unit]
Description=Send Discord webhook on system boot
After=network-online.target tailscaled.service time-sync.target
Wants=network-online.target tailscaled.service time-sync.target

[Service]
Type=oneshot
User=%i
ExecStart=/home/%i/bin/discord_startup.sh
Environment=PATH=/usr/bin:/bin:/usr/local/bin

[Install]
WantedBy=multi-user.target

